# sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :0 -screen 0 1280x720x24 -nolisten tcp -ac
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :0 -forever -shared -localhost -rfbport 5900 -nopw
priority=20
startsecs=2
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify]
command=websockify --web /usr/share/novnc 127.0.0.1:8080 127.0.0.1:5900
priority=30
startsecs=2
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:chromium]
command=/usr/local/bin/dev-browser-chromium --user-data-dir=/data --no-first-run --disable-session-crashed-bubble --window-position=0,0 --window-size=1280,720 --remote-debugging-address=127.0.0.1 --remote-debugging-port=9222 --no-sandbox --disable-dev-shm-usage
priority=40
startsecs=3
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sandbox-proxy]
command=/usr/local/bin/sandbox-proxy
priority=50
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
